# Firmware/U2W_Update_28102020/etc/inittab
console::sysinit:/etc/init.d/rcS
#ttyS0::respawn:/sbin/getty -n -l/bin/sh -L ttyS0 115200 vt100 # GENERIC_SERIAL
::askfirst:-/bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
#::shutdown:bin/umount -a -r
::shutdown:/script/after_shutdown.sh

